# systemd/presence-monitor.service
[Unit]
Description=Garage Presence Monitor
After=network.target

[Service]
Type=simple
User=tim
WorkingDirectory=/home/tim/projects/20260119_is_the_garage_door_open/is_the_garage_door_open
Environment="PYTHONUNBUFFERED=1"
StandardOutput=journal
StandardError=journal
ExecStart=/usr/bin/python3 -u presence_monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
